# install/com.dwcts.pixel-agents.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dwcts.pixel-agents</string>

    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__REPO_DIR__/standalone/bin/cli.mjs</string>
        <string>--no-open</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__/standalone</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PIXEL_AGENTS_PORT</key>
        <string>__PORT__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/pixel-agents.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/pixel-agents.err.log</string>
</dict>
</plist>
